# wg-failover.service
[Unit]
Description=WireGuard Failover Service
Documentation=https://github.com/yourusername/wg-failover
After=network-online.target NetworkManager.service systemd-networkd.service
Wants=network-online.target
ConditionPathExists=/usr/local/bin/wg-failover

[Service]
Type=simple
ExecStart=/usr/local/bin/wg-failover
Environment="WG_FAILOVER_CONFIG=/etc/wg-failover/config.toml"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
User=root
Group=root

# Security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

# The service uses the configuration file at /etc/wg-failover/config.toml
# Make sure to create and configure this file before starting the service

[Install]
WantedBy=multi-user.target
Alias=wireguard-failover.service